# airis-mcp-gateway/airis-mcp-gateway.service
[Unit]
Description=Airis MCP Gateway
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=tovic
ExecStart=/home/tovic/.local/bin/airis-start.sh

[Install]
WantedBy=multi-user.target
